# laravel-mqtt-subscriber-local-ordermac.conf
[program:subscribe-local-ordermac]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan mqtt:subscribe-local-ordermac
autostart=true
autorestart=true
startretries=3
user=www-data
stderr_logfile=/var/www/html/storage/logs/subscribe-local-ordermac.err.log
stdout_logfile=/var/www/html/storage/logs/subscribe-local-ordermac.out.log
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
