# laravel-publish-order-stats.conf
[program:publish-order-stats]
process_name=%(program_name)s_%(process_num)02d
#command=/var/www/html/scripts/modbus-loop.sh
command=php /var/www/html/artisan mqtt:publish-order-stats
autostart=true
autorestart=true
startretries=0
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/laravel-publish-order-stats.out.log
stderr_logfile=/var/www/html/storage/logs/laravel-publish-order-stats.err.log
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
