# laravel-232-basculas-rs232.conf
[program:232-basculas-rs232]
process_name=%(program_name)s_%(process_num)02d
command=node /var/www/html/232-basculas/index.js
directory=/var/www/html/232-basculas
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/232-basculas-rs232.log
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10
priority=100
startsecs=5
startretries=3
