# laravel-modbus-web-8001.conf
[program:modbus-web.8001]
process_name=%(program_name)s_%(process_num)02d
command=python3 /var/www/html/modbus-web-8001.py
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/modbus-web.log
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10